# deploy/systemd/codex-feishu-bridge.service.in
[Unit]
Description=Codex Feishu Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@WORKDIR@
EnvironmentFile=%h/.config/codex-feishu-bridge/bridge.env
Environment=HOME=@HOME@
Environment=TMPDIR=/tmp
Environment=PATH=@PATH@
ExecStart=@NODE_BIN@ --enable-source-maps @WORKDIR@/dist/index.js
TimeoutStopSec=15
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
