# deploy/systemd/lark-codex.service.in
[Unit]
Description=Lark Codex
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@HOME@
EnvironmentFile=%h/.config/lark-codex/bridge.env
Environment=HOME=@HOME@
Environment=TMPDIR=/tmp
Environment=PATH=@PATH@
ExecStart=@BIN_PATH@
TimeoutStopSec=15
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
