# deploy/launchd/com.lark-codex.bridge.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lark-codex.bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>@RUNNER_PATH@</string>
  </array>

  <key>WorkingDirectory</key>
  <string>@WORKDIR@</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>@STDOUT_PATH@</string>

  <key>StandardErrorPath</key>
  <string>@STDERR_PATH@</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>@HOME@</string>
    <key>PATH</key>
    <string>@PATH@</string>
    <key>TMPDIR</key>
    <string>/tmp</string>
  </dict>
</dict>
</plist>
